# pi-server.service
[Unit]
Description=pi-server OpenAI-compatible API for pi-coding-agent
Documentation=https://github.com/ldarren/pi-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claw01
Group=claw01
WorkingDirectory=/home/claw01/mygit/pi-server
Environment=NODE_ENV=production
Environment=PATH=/home/claw01/.nvm/versions/node/v22.22.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/claw01/.nvm/versions/node/v22.22.1/bin/npm start
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
